# tangelo/assets/daemon/systemd/system/tangelo@.service
[Unit]
Description=Tangelo Web Server on %I
After=network.target

[Service]
ExecStart=/usr/lib/systemd/scripts/launch-tangelo.sh %I
